# packages/DEBIAN/prerm
#!/bin/sh
set -e

echo "Pre-Removal Macro"

if [ "$1" = "remove" ]; then
	systemctl disable __PACKAGE_NAME__-worker || true
	systemctl stop __PACKAGE_NAME__-worker || true

	systemctl disable __PACKAGE_NAME__ || true
	systemctl stop __PACKAGE_NAME__ || true

	rm -f __PACKAGE_DIR__/config/settings.cfg
fi

exit 0
